# api/supervisor-app.conf
[program:app-gunicorn]
command = gunicorn -c gunicorn.settings.py weather.wsgi
directory = /app

[program:nginx-app]
command = /usr/sbin/nginx
